# dist/aa-cream_1.0-1/lib/systemd/system/creamd.service
[Unit]
Description=Constant Recording Easy Access to Media
After=network.target

[Service]
ExecStart=/opt/austrianAudio/bin/cream start
ExecReload=/opt/austrianAudio/bin/cream restart
ExecStop=/opt/austrianAudio/bin/cream stop
KillMode=control-group
Restart=on-failure
RestartPreventExitStatus=255
Type=forking
CacheDirectory=CREAM
User=root

[Install]
WantedBy=multi-user.target
Alias=creamd.service

